# xtc/firstboot-late.service
[Unit]
Description=FirstBoot-late
After=multi-user.target
ConditionFileNotEmpty=/opt/thinclient/install/firstboot-late.sh

[Service]
ExecStart=/opt/thinclient/install/firstboot-late.sh
ExecStartPost=/bin/mv /opt/thinclient/install/firstboot-late.sh /opt/thinclient/install/firstboot-late.sh.done
Type=oneshot
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
